# Frontend/setup/cmpe281_frontend.service
[Unit]
Description=CMPE281-CloudDrive-Frontend
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=bash /home/ubuntu/Frontend/run.sh

[Install]
WantedBy=multi-user.target